# units/llama-gemma4.service
[Unit]
Description=llama.cpp ROCm seat: Gemma 4-26B QAT (q4_0, GPU)
After=network.target

[Service]
Type=simple
User=llama
Group=llama-models
SupplementaryGroups=render video
UMask=002
WorkingDirectory=/opt/llama
EnvironmentFile=-/etc/llama-fleet.env
Environment="ROCM_HOME=/opt/rocm-7.2.2"
Environment="PYTORCH_ROCM_ARCH=gfx1151"
Environment="HSA_OVERRIDE_GFX_VERSION=11.5.1"
Environment="AMD_SERIALIZE_KERNEL=3"
ExecStart=/opt/llama/bin/llama-server-rocm \
  -m /opt/models/gemma-4-26B_q4_0-it.gguf \
  -ngl all \
  -c 32768 \
  --port 8013 \
  --host 127.0.0.1 \
  --threads 16 \
  --threads-batch 16 \
  --parallel 2 \
  --cont-batching \
  --metrics \
  --cache-ram 512 \
  --cache-type-k q8_0 \
  --cache-type-v q8_0 \
  --flash-attn on \
  --slot-save-path /var/lib/llama/hotcache/gemma4/ \
  --jinja
Restart=on-failure
RestartSec=15
# llama-server ignores SIGTERM under load; this is the wall before SIGKILL.
TimeoutStopSec=120

MemoryMax=16G
MemoryLow=4G

[Install]
WantedBy=multi-user.target
